# src/watchdog/watchdog.service
[Unit]
Description=VPN Watchdog Service - Bandwidth Monitor
After=network.target sing-box.service

[Service]
Type=simple
User=root
WorkingDirectory=/home/ubuntu/vpn-bot
Environment="PYTHONUNBUFFERED=1"
ExecStart=/usr/bin/python3 /home/ubuntu/vpn-bot/watchdog/service.py
Restart=always
RestartSec=60

[Install]
WantedBy=multi-user.target
